# .config/systemd/user/slack.service
[Unit]
Description=Slack

[Service]
Environment=XAUTHORITY=%h/.Xauthority
Environment=DISPLAY=:0
ExecStart=/usr/bin/slack
CPUQuotaPeriodSec=10ms
